# services/flask_updater.conf
description "flask updater"
start on stopped rc RUNLEVEL=[2345]
respawn
exec python /home/[username]/[path]/webhook/flask_updater.py